# .emergent/applyfixes.conf
# /app/.emergent/applyfixes.conf
# Supervisor program that runs apply.sh once at boot, then exits.
# This file is copied to /etc/supervisor/conf.d/ by apply.sh or the boot hook.
[program:applyfixes]
command=/bin/bash /app/.emergent/apply.sh
directory=/app
autostart=true
autorestart=false
startsecs=0
startretries=0
exitcodes=0
priority=1
stderr_logfile=/var/log/supervisor/applyfixes.err.log
stdout_logfile=/var/log/supervisor/applyfixes.out.log
